# systemd/chrome-borderless.service
[Unit]
Description=Chrome Borderless HTTP Server
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/%u/.config/chrome-borderless/daemon-wrapper.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
