# systemd/chrome-borderless.service
[Unit]
Description=Chrome Borderless HTTP Server

[Service]
Type=simple
ExecStart=/home/%u/.config/chrome-borderless/daemon.py
Restart=on-failure
RestartSec=5
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="DISPLAY=:0"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
